# deploy/systemd/dharmagptbeta-server.service.template
[Unit]
Description=DharmaGPT beta FastAPI server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={{REPO_ROOT}}/dharmagpt
EnvironmentFile={{REPO_ROOT}}/dharmagpt/.env
ExecStart={{REPO_ROOT}}/dharmagpt/.venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
